# vdrpbd/vdrpbd.service
[Unit]
Description=VDR ACPI power button handling daemon

[Service]
ExecStart=/usr/bin/vdrpbd -f

[Install]
WantedBy=multi-user.target
