# vdradmin-am/vdradmin-am.service
[Unit]
Description=VDRAdmin-AM
After=vdr.service

[Service]
Type=forking
ExecStart=/usr/bin/vdradmind

[Install]
WantedBy=multi-user.target