# epgd/epgd.service
[Unit]
Description=VDR EPG Daemon
Wants=mysqld.service

[Service]
ExecStart=/usr/bin/epgd

[Install]
WantedBy=multi-user.target